# crontab
* */6 * * * ./podcast-v3/podcast.sh >/dev/null 2>&1
